# accounting-backend.service
[Unit]
Description=Almasa Accounting Backend API
Documentation=https://github.com/yourusername/accounting-app
After=network.target mysql.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/var/www/accounting-app/backend
Environment=NODE_ENV=production
EnvironmentFile=/var/www/accounting-app/backend/.env
ExecStart=/usr/bin/pm2 start ecosystem.config.js --env production
ExecReload=/usr/bin/pm2 reload accounting-backend
ExecStop=/usr/bin/pm2 stop accounting-backend
Restart=on-failure
RestartSec=10s

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=accounting-backend

[Install]
WantedBy=multi-user.target
